# roles/topbeat/files/init_topbeat.conf
# Topbeat
start on runlevel [2345] 
stop on runlevel [!2345] 

script
    exec /root/shipper/topbeat/topbeat -c /root/shipper/topbeat/topbeat.yml
end script
